# prosody.service
[Unit]
Description=Prosody Speech-to-Text Service
After=graphical-session.target

[Service]
Type=simple
# Note: This is a template. The install.sh script will create a user-specific version
# with the correct path to your prosody installation
ExecStart=%h/.local/bin/prosody
Restart=always
RestartSec=10
Environment="DISPLAY=:0"

[Install]
WantedBy=default.target