# debian/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        # Copy service file to system location
        cp /usr/share/prosody-stt/prosody.service /etc/systemd/user/prosody.service || true
        
        echo ""
        echo "✓ Prosody has been installed!"
        echo ""
        echo "To start using Prosody:"
        echo "  1. Run: prosody"
        echo "  2. Enable auto-start: systemctl --user enable prosody"
        echo ""
        echo "Usage:"
        echo "  - Double-tap Ctrl to start/stop recording"
        echo "  - Double-tap Escape to cancel recording"
        echo ""
        ;;
esac

#DEBHELPER#

exit 0