# debian/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        # Copy service file to system location
        cp /usr/share/prosody-stt/prosody.service /etc/systemd/user/prosody.service || true
        
        # Enable auto-start for the user who installed it
        if [ -n "$SUDO_USER" ]; then
            su - "$SUDO_USER" -c "systemctl --user daemon-reload 2>/dev/null || true"
            su - "$SUDO_USER" -c "systemctl --user enable prosody.service 2>/dev/null || true"
        fi
        
        echo ""
        echo "✓ Prosody has been installed!"
        echo ""
        echo "To start using Prosody: run 'prosody'"
        echo "Auto-start on login has been enabled."
        echo ""
        echo "Usage:"
        echo "  - Double-tap Ctrl to start/stop recording"
        echo "  - Double-tap Escape to cancel recording"
        echo ""
        ;;
esac

#DEBHELPER#

exit 0