# debian/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        # Create system user for prosody if needed (optional)
        # For now, we run as the regular user
        
        # Ensure python3-pip is available
        if ! command -v pip3 &> /dev/null; then
            echo "Installing pip3..."
            apt-get update
            apt-get install -y python3-pip
        fi
        
        # Install Python package in system
        pip3 install --system /usr/share/prosody-stt/ || true
        
        # Enable auto-start for all users
        echo "Enabling auto-start for Prosody..."
        
        # Create systemd user service directory if it doesn't exist
        mkdir -p /etc/systemd/user/
        
        # Copy service file to system location
        cp /usr/share/prosody-stt/prosody.service /etc/systemd/user/prosody.service || true
        
        # For the user running the install, enable and start the service
        if [ -n "$SUDO_USER" ]; then
            # Running under sudo, enable for the actual user
            su - "$SUDO_USER" -c "systemctl --user daemon-reload" || true
            su - "$SUDO_USER" -c "systemctl --user enable prosody.service" || true
            su - "$SUDO_USER" -c "systemctl --user start prosody.service" || true
        else
            # Not under sudo, enable for current user
            systemctl --user daemon-reload || true
            systemctl --user enable prosody.service || true
            systemctl --user start prosody.service || true
        fi
        
        echo ""
        echo "✓ Prosody has been installed and is now running!"
        echo ""
        echo "Usage:"
        echo "  - Double-tap Ctrl to start/stop recording"
        echo "  - Double-tap Escape to cancel recording"
        echo ""
        echo "To manage Prosody:"
        echo "  - Check status: systemctl --user status prosody"
        echo "  - Stop: systemctl --user stop prosody"
        echo "  - Disable auto-start: systemctl --user disable prosody"
        echo ""
        echo "Note: First run will download the Whisper model (~140MB)"
        ;;
esac

#DEBHELPER#

exit 0